# part4/moj-debian/config/includes.chroot/etc/systemd/system/update-issue.service
[Unit]
Description=Update issue with current IP
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-issue.sh

[Install]
WantedBy=multi-user.target
